# frontends/server-signal-helper/esdm-server-suspend.service.in
# Systemd configuration file
#
# Notify the ESDM server about suspend/sleep operation
#

[Unit]
Description=Entropy Source and DRNG Manager Daemon Suspend/Sleep Trigger
After=sleep.target
After=suspend.target
After=hibernate.target
After=esdm-server.service
Wants=esdm-server.service
ConditionPathExists=/run/esdm-server.pid

[Service]
Type=oneshot
ExecStart=@PATH@/esdm-server-signal-helper --suspend --pid /run/esdm-server.pid
TimeoutSec=0

# Filesystem
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ReadOnlyPaths=-/

# Network
IPAddressDeny=any
PrivateNetwork=yes
RestrictAddressFamilies=AF_UNIX

# Process isolation
ProtectProc=invisible
ProcSubset=pid
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege hardening
NoNewPrivileges=yes
RestrictSUIDSGID=yes
CapabilityBoundingSet=CAP_KILL
AmbientCapabilities=CAP_KILL
LimitMEMLOCK=0

# Syscall filtering
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@aio @chown @clock @cpu-emulation @debug @ipc @keyring @memlock @module @mount @obsolete @pkey @privileged @raw-io @reboot @resources @sandbox @setuid @swap

# Misc hardening
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RemoveIPC=yes
UMask=0077

[Install]
WantedBy=sleep.target
WantedBy=suspend.target
WantedBy=hibernate.target
