# frontends/cuse/esdm-proc.service.in
# Systemd configuration file
#
# The ESDM /proc/sys/kernel/random provider does not depend on any other system
# service except the ESDM server. Furthermore, we want to start it as early as
# possible in the boot cycle so that other services requiring the proc file
# system (e.g. entropy gathering daemons) can rely on it.

[Unit]
Description=Entropy Source and DRNG Manager /proc/sys/kernel/random
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target
StartLimitIntervalSec=500
StartLimitBurst=5
Wants=esdm-server.service
After=esdm-server.service
PartOf=esdm-linux-compat.target

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=@PATH@/esdm-proc --relabel -f -o allow_other /proc/sys/kernel/random
KillMode=mixed
