# frontends/cuse/esdm-cuse-random.service.in
# Systemd configuration file
#
# The ESDM /dev/random provider does not depend on any other system
# service except the ESDM server. Furthermore, we want to start it as early as
# possible in the boot cycle so that other services requiring
# random numbers (like SSHD, TLS-based services) benefit from it

[Unit]
Description=Entropy Source and DRNG Manager /dev/random
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target
StartLimitIntervalSec=500
StartLimitBurst=5
Wants=esdm-server.service
After=esdm-server.service
PartOf=esdm-linux-compat.target

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=@PATH@/esdm-cuse-random -f
KillMode=mixed
